# investment_today_intent/supervisord.conf
[supervisord]
nodaemon=true

[program:intent]
command=python3 api.py
autostart=true
autorestart=true
stdout_logfile = /var/log/stdout/intent.log
redirect_stderr=true


